# Library/LaunchDaemons/local.all.dnsmasq.restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>homebrew.mxcl.dnsmasq.restart</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/launchctl</string>
      <string>kill</string>
      <string>SIGHUP</string>
      <string>system/homebrew.mxcl.dnsmasq</string>
    </array>
    <key>WatchPaths</key>
    <array>
      <string>/Users/camdennarzt/Developer/Bash/dot-files/usr/local/etc/dnsmasq.d</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
  </dict>
</plist>
