# Library/LaunchDaemons/local.secur-t.dyn_ip.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.secur-t.dns</string>
    <key>UserName</key>
    <string>camdennarzt</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/dyn_ip</string>
    </array>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/log/dyn_ip.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/log/dyn_ip.log</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>
  </dict>
</plist>
