# Library/LaunchAgents/com.arewesmallyet.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.arewesmallyet.update</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-i</string>
      <string>-c</string>
      <string>$SHELL --login -i -c "bundle exec padrino rake update_once"</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/camdennarzt/Developer/Ruby/arewesmallyet</string>
  </dict>
</plist>
