# Library/LaunchAgents/com.arewesmallyet.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.arewesmallyet.update</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/camdennarzt/.rbenv/shims/ruby</string>
      <string>/Users/camdennarzt/Developer/Ruby/firefox.rb</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/camdennarzt/Developer/Ruby/arewesmallyet</string>
    <key>StandardOutPath</key>
    <string>/Users/camdennarzt/Developer/Ruby/arewesmallyet/log/update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/camdennarzt/Developer/Ruby/arewesmallyet/log/update_err.log</string>
  </dict>
</plist>
