# Library/LaunchAgents/local.currency.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>currency.update</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/camdennarzt/Developer/Bash/dot-files/Library/LaunchAgents/currencyUpdate.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/log/CurrencyUpdate/log.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/CurrencyUpdate/error.log</string>
  </dict>
</plist>
