# Library/LaunchAgents/local.secur-t.speedtest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.secur-t.speedtest</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/camdennarzt/Developer/Ruby/speedtest.rb</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/camdennarzt/Library/Logs/speedtest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/camdennarzt/Library/Logs/speedtest_err.log</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>
  </dict>
</plist>
